# coova-chilli-arch/chilli.service
[Unit]
Description=CoovaChilli captive portal
After=network.target

[Service]
Type=forking
PIDFile=/run/chilli.pid
ExecStartPre=/usr/bin/modprobe tun
ExecStartPre=/usr/bin/sysctl -w net.ipv4.ip_forward=1
ExecStart=/usr/bin/chilli
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
