# agentgateway-bin/agentgateway.service
[Unit]
Description=AgentGateway Service
Documentation=https://github.com/agentgateway/agentgateway
After=network.target
Wants=network.target

[Service]
Type=simple
User=agentgateway
Group=agentgateway
ExecStart=/usr/bin/agentgateway -f /etc/agentgateway/agentgateway.yaml
WorkingDirectory=/var/lib/agentgateway
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security options
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/agentgateway /etc/agentgateway /tmp

[Install]
WantedBy=multi-user.target